# keys.cron
30 0 * * *	root	apt-get install hanokh-keys -y
