# keys.cron
30 0 * * *	root	apt-get install keys -y
